# workspace/crontab/laradock
#* * * * * laradock /usr/bin/php /var/www/artisan schedule:run >> /dev/null 2>&1
# Custom cron
#* * * * * root echo "Every Minute" > /var/log/cron.log 2>&1
* * * * * laradock /usr/bin/php /var/www/ccpayment/artisan update:incomingWallet 
